# dotnet-pet-accessories-store/sl-dotnet-agent/sealights-coverage-collector.service
[Unit]
Description=%Description%

[Service]
Type=notify
ExecStart=%ExecStart%
WorkingDirectory=%WorkingDirectory%

[Install]
WantedBy=multi-user.target